# vllm_monitor/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
priority=1

[program:loki]
command=/usr/local/bin/loki -config.file=/etc/loki/loki-config.yml
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/loki-stdout.log
stderr_logfile=/var/log/supervisor/loki-stderr.log
priority=2

[program:promtail]
command=/usr/local/bin/promtail -config.file=/etc/promtail/promtail-config.yml
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/promtail-stdout.log
stderr_logfile=/var/log/supervisor/promtail-stderr.log
priority=3
depends_on=loki

[program:vllm]
command=/usr/local/bin/start-vllm.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/vllm/vllm.log
stderr_logfile=/var/log/vllm/vllm-error.log
priority=4
